# smartctl_exporter/rootfs/etc/services.d/smartctl_exporter/run
#!/usr/bin/env bashio
# ==============================================================================
# Home Assistant Community Add-on: smartctl_exporter
# Runs smartctl_exporter
# ==============================================================================

# Set bashio log level from add-on configuration
bashio::log.level "$(bashio::config 'log_level')"

bashio::log.info "Starting smartctl_exporter..."

# Build command flags array
declare -a flags=(
    --smartctl.interval="$(bashio::config 'smartctl_interval')"
    --smartctl.rescan="$(bashio::config 'smartctl_rescan')"
    --log.level="$(bashio::config 'log_level')"
    --web.listen-address=":9633"
)

bashio::log.debug "Configured interval: $(bashio::config 'smartctl_interval')"
bashio::log.debug "Configured rescan: $(bashio::config 'smartctl_rescan')"
bashio::log.debug "Configured log level: $(bashio::config 'log_level')"

# Add device exclusion pattern if configured
if bashio::config.has_value 'smartctl_device_exclude'; then
    device_exclude=$(bashio::config 'smartctl_device_exclude')
    bashio::log.debug "Excluding devices matching pattern: ${device_exclude}"
    flags+=(--smartctl.device-exclude="${device_exclude}")
else
    bashio::log.debug "No device exclusion pattern configured"
fi

# Add device inclusion pattern if configured
if bashio::config.has_value 'smartctl_device_include'; then
    device_include=$(bashio::config 'smartctl_device_include')
    bashio::log.debug "Including only devices matching pattern: ${device_include}"
    flags+=(--smartctl.device-include="${device_include}")
else
    bashio::log.debug "No device inclusion pattern configured"
fi

# Add specific devices if configured
if bashio::config.has_value 'smartctl_devices'; then
    bashio::log.debug "Processing smartctl_devices configuration"
    devices_list=$(bashio::config 'smartctl_devices')
    bashio::log.debug "Devices list:\n${devices_list}"
    device_count=0
    while read -r device; do
        if bashio::var.has_value "${device}"; then
            bashio::log.debug "Adding device: ${device}"
            flags+=(--smartctl.device="${device}")
            device_count=$((device_count + 1))
        fi
    done <<< "${devices_list}"
else
    bashio::log.debug "No specific devices configured"
fi

bashio::log.info "Starting smartctl_exporter with flags: ${flags[*]}"

# Run smartctl_exporter
exec /usr/local/bin/smartctl_exporter "${flags[@]}"
