# smartctl_exporter/rootfs/etc/services.d/smartctl_exporter/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: smartctl_exporter
# Runs smartctl_exporter
# ==============================================================================

bashio::log.info "Starting smartctl_exporter..."

# Build command flags array
declare -a flags=(
    --smartctl.interval="$(bashio::config 'smartctl_interval')"
    --smartctl.rescan="$(bashio::config 'smartctl_rescan')"
    --log.level="$(bashio::config 'log_level')"
    --web.listen-address=":9633"
)

# Add device exclusion pattern if configured
if bashio::config.has_value 'smartctl_device_exclude'; then
    flags+=(--smartctl.device-exclude="$(bashio::config 'smartctl_device_exclude')")
fi

# Add device inclusion pattern if configured
if bashio::config.has_value 'smartctl_device_include'; then
    flags+=(--smartctl.device-include="$(bashio::config 'smartctl_device_include')")
fi

# Add specific devices if configured
if bashio::config.has_value 'smartctl_devices'; then
    while read -r device; do
        if bashio::var.has_value "${device}"; then
            bashio::log.info "Adding device: ${device}"
            flags+=(--smartctl.device="${device}")
        fi
    done < <(bashio::config 'smartctl_devices | .[]')
fi

bashio::log.info "Starting smartctl_exporter with flags: ${flags[*]}"

# Run smartctl_exporter
exec /usr/local/bin/smartctl_exporter "${flags[@]}"
